# deploy/habitat-api.service
[Unit]
Description=Habitat API

[Service]
Type=simple
WorkingDirectory=%h/habitat-cli
ExecStart=/usr/bin/env HABITAT_API_HOST=0.0.0.0 HABITAT_API_PORT=8787 /usr/local/bin/bun run server
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
